# infra/systemd/ibgateway.service
[Unit]
Description=IB Gateway (Paper Trading) via IBC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment="DISPLAY=:1"
ExecStartPre=-/usr/bin/pkill -f 'Xvfb :1'
ExecStartPre=/bin/bash -c '/usr/bin/Xvfb :1 -screen 0 1024x768x24 &'
ExecStart=/opt/ibc/gatewaystart.sh -inline
ExecStopPost=-/usr/bin/pkill -f 'Xvfb :1'
KillMode=control-group
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
