# infra/systemd/trading-relay.service
[Unit]
Description=Trading Assistant Relay
After=network.target

[Service]
Type=simple
User=trading-relay
Group=trading-relay
WorkingDirectory=/opt/trading-relay
ExecStart=/opt/trading-relay/start.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/trading-relay/data
PrivateTmp=true

# Environment
EnvironmentFile=-/opt/trading-relay/.env

[Install]
WantedBy=multi-user.target
